# docker/supervisor.d/cronjob.conf
[program:cronjob]
command=/usr/sbin/crond -f -l 8
numprocs=1
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
